# shoot.service
# exists as /etc/systemd/system/shoot.service

[Unit]
Description=Take pictures a lot

[Service]
ExecStart=/home/pi/still.sh

